# cron/second-brain.cron
# personal-brain backup
0 * * * * /bin/bash /home/lalitmee/dotfiles/scripts/backup/run-cron-jobs personal >/dev/null 2>&1

# work-brain backup
0 * * * * /bin/bash /home/lalitmee/dotfiles/scripts/backup/run-cron-jobs work >/dev/null 2>&1

# * * * * * /bin/bash /home/lalitmee/dotfiles/scripts/test/test-home

# * * * * * /bin/bash /home/lalitmee/dotfiles/scripts/backup/test-cron
