# cron/second-brain.cron
0 * * * * /bin/bash /Users/lalit.kumar1/dotfiles/scripts/backup/run-cron-jobs >> /Users/lalit.kumar1/dotfiles/scripts/logs/second-brain.log 2>&1
